# com.openclaw.octen-embeddings.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.octen-embeddings</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ms/.local/share/mise/shims/python3</string>
        <string>/Users/ms/personal/octen-embeddings-server/server.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ms/personal/octen-embeddings-server</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>MLX_MODEL_PATH</key>
        <string>/Users/ms/personal/octen-embeddings-server/models/Octen-Embedding-8B-mlx</string>
        <key>PORT</key>
        <string>8100</string>
        <key>HOST</key>
        <string>127.0.0.1</string>
        <key>LOG_LEVEL</key>
        <string>INFO</string>
        <key>PATH</key>
        <string>/Users/ms/.local/share/mise/shims:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ms/Library/Logs/octen-embeddings.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ms/Library/Logs/octen-embeddings.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
